# templates/upstart.conf
#!upstart
description "Apache Zeppelin Server"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on shutdown

respawn

pre-start script
su ubuntu -c "cd {{zeppelin_home}}; {{zeppelin_home}}/bin/zeppelin-daemon.sh --config {{zeppelin_conf}} start"
end script

post-stop script
su ubuntu -c "{{zeppelin_home}}/bin/zeppelin-daemon.sh --config {{zeppelin_conf}} stop"
end script
